# deploy/systemd/ai-tech-daily-brief-push-retry.service
[Unit]
Description=Retry AI tech daily brief PDF push after failure
Wants=network-online.target docker.service
After=network-online.target docker.service
OnFailure=ai-tech-daily-brief-failure@%n.service

[Service]
Type=oneshot
EnvironmentFile=/etc/ai-tech-daily-brief/env
ExecStartPre=/usr/bin/test -f /etc/ai-tech-daily-brief/env
ExecStartPre=/usr/bin/docker image inspect ${AI_TECH_DAILY_BRIEF_IMAGE}
ExecCondition=/bin/sh -c 'test -f /var/lib/ai-tech-daily-brief/state/$(date +%%F).push.failed'
ExecCondition=/bin/sh -c '! test -f /var/lib/ai-tech-daily-brief/state/$(date +%%F).push.sent'
ExecCondition=/bin/sh -c 'test -s /var/lib/ai-tech-daily-brief/output/ai-tech-daily-brief-$(date +%%F).pdf'
ExecStart=/usr/bin/docker run --rm --name ai-tech-daily-brief-push-retry --env-file /etc/ai-tech-daily-brief/env -e OUTPUT_ROOT=/app/output -e LOG_ROOT=/app/logs -e TMP_ROOT=/app/tmp -e STATE_ROOT=/app/state -v /var/lib/ai-tech-daily-brief/output:/app/output:ro -v /var/log/ai-tech-daily-brief:/app/logs -v /var/lib/ai-tech-daily-brief/state:/app/state ${AI_TECH_DAILY_BRIEF_IMAGE} push
StandardOutput=append:/var/log/ai-tech-daily-brief/push-retry.log
StandardError=append:/var/log/ai-tech-daily-brief/push-retry.err.log
